# misc/faxgetty-ttyT38-10.service
[Unit]
Description=HylaFAX faxgetty for ttyT38-10

[Service]
User=root
Group=root
Restart=always
RestartSec=30
ExecStart=/usr/sbin/faxgetty ttyT38-10

[Install]
WantedBy=multi-user.target
